# com.user.copilot-update-checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.copilot-update-checker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>~/.local/bin/copilot_update_checker.scpt</string>
    </array>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/copilot-update-checker.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/copilot-update-checker.error.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>